# systemd/glpi-ad-integration.service
[Unit]
Description=Integração GLPI - Active Directory (Criação de Usuários)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=glpi-ad-svc
Group=glpi-ad-svc
WorkingDirectory=/opt/glpi-ad-integration
EnvironmentFile=/opt/glpi-ad-integration/.env
ExecStart=/opt/glpi-ad-integration/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8443 --workers 2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=glpi-ad-integration

# Segurança
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/glpi-ad-integration
ReadOnlyPaths=/etc/ssl/certs /etc/ssl/private /usr/share/ca-certificates
PrivateTmp=true

[Install]
WantedBy=multi-user.target
